# scripts/companion.service
[Unit]
Description=Atlas AI Companion — Voice-driven AI face on Raspberry Pi 5
After=network-online.target sound.target bluetooth.service
Wants=network-online.target

[Service]
Type=simple
User=daveealex
Group=daveealex
WorkingDirectory=/home/daveealex/ai-companion
Environment=PATH=/home/daveealex/ai-companion/venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=WEB_SERVER_URL=http://localhost:5000
Environment=DISPLAY=:0

# Start the web server, then chromium kiosk, then orchestrator
ExecStartPre=/bin/sleep 2
ExecStart=/bin/bash -c 'cd /home/daveealex/ai-companion && source venv/bin/activate && python3 web/server.py &'
ExecStartPost=/bin/sleep 2
ExecStart=/bin/bash -c 'cd /home/daveealex/ai-companion && chromium-browser --kiosk --no-first-run --disable-features=TranslateUI http://localhost:5000 &'
ExecStartPost=/bin/sleep 3
ExecStart=/bin/bash -c 'cd /home/daveealex/ai-companion && source venv/bin/activate && python3 orchestrator.py'

Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=atlas-companion

# Need audio and display access
AmbientCapabilities=CAP_AUDIO
# Allow Chromium to run headless
Environment=CHROMIUM_FLAGS=--kiosk --no-first-run --disable-gpu --disable-features=TranslateUI,Notifications

[Install]
WantedBy=multi-user.target
