# scripts/companion.service
[Unit]
Description=Red Queen AI Companion
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
User=daveealex
Group=daveealex
WorkingDirectory=/home/daveealex/ai-companion
Environment=PATH=/home/daveealex/ai-companion/venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=WEB_SERVER_URL=http://localhost:5000
Environment=DISPLAY=:0
Environment=CHROMIUM_FLAGS=--kiosk --no-first-run --disable-gpu --disable-features=TranslateUI,Notifications

ExecStartPre=/bin/sleep 2
ExecStart=/home/daveealex/ai-companion/scripts/launch.sh

Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=red-queen-companion

AmbientCapabilities=CAP_SYS_NICE

[Install]
WantedBy=multi-user.target
